# app-clip/hotplug-clip/files/_debian/prerm
#!/bin/sh

/sbin/rc-update del mkdevices default || echo "! rc-upadte del mkdevices failed !" >&2
