# app-clip/clip-netd/files/netup
#!/sbin/runscript

depend() {
	need ipsec
}

start() {
	einfo "Networking services started"
}
