# app-clip/clip-useradmin/files/_debian/postinst
#!/bin/sh

/sbin/rc-update -a useradmin default

