# app-clip/clip-useradmin/files/_debian/prerm
#!/bin/sh

/sbin/rc-update -d useradmin
