# app-clip/clip-mdadm/files/_debian/postinst
#!/bin/sh

/sbin/rc-update -a mdadmd default

